# job/crontab
* * * * * sh /script/test.sh
@reboot sh /script/test.sh
